# cron/crontab
*/10 * * * * sudo chown plante: /dev/ttyACM0; /usr/bin/python /home/plante/tt_hr_sensor_network/python/serial_read_from_arduino.py

# To debug add > /home/plante/log.txt 2>&1
